# package/alcatel-dongle/alcatel-dongle
#!/sbin/openrc-run

description="Daemon that creates reverse-proxy on wlan0 for web interface \
of alcatel dongle"

command_background=true

depend() {
  use logger
}

start() {
  ebegin "Starting alcatel dongle"
  einfo "Bringing wlan0 up"
  ip link set dev wlan0 up
  einfo "Setting subnet"
  ip addr add 10.0.0.1/24 dev wlan0

  einfo "Starting dnsmasq"
  dnsmasq --conf-file=/etc/dnsmasq.conf

  einfo "Starting hostapd"
  hostapd -B -P /var/run/hostapd.pid /etc/hostapd.conf

  einfo "Starting tinyproxy"
  tinyproxy -c /etc/tinyproxy/tinyproxy.conf
}

stop() {
  kill $(cat /var/run/dnsmasq.pid)
  kill $(cat /var/run/hostapd.pid)
  kill $(cat /var/run/tinyproxy.pid)
}
