# board/mangopi/rootfs-overlay/backup/etc/cron.daily/docker-backup.sh
#!/bin/sh

mkdir -p /backup
tar -c /docker | zstd > /backup/$(date -Ihours).tar.zstd
