# board/mangopi/rootfs-overlay/backup/etc/cron.daily/docker-backup.sh
#!/bin/sh

cp -r /docker /backup/$(date --iso-8601=hours)
